# deploy/scaleway/systemd/visual-compare-api.service
[Unit]
Description=visual-compare API (Express + SQLite + Playwright + ImageMagick)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=visual-compare
Group=visual-compare
WorkingDirectory=/opt/visual-compare/packages/api

# Secrets + config live in this file. Deploy step writes/refreshes it.
# Format: KEY=value per line; no quoting (systemd EnvironmentFile semantics).
EnvironmentFile=/etc/visual-compare/env

# These two are intentionally hard-coded (not in the env file) so an
# operator can't accidentally point the prod API at a dev DB.
Environment=DB_PATH=/mnt/data/visual-compare.sqlite
Environment=IMAGES_DIR=/mnt/data/images
Environment=LM_LAST_USE_PATH=/mnt/data/lm-last-use
Environment=LISTEN_HOST=127.0.0.1
Environment=NODE_ENV=production
# Match the PLAYWRIGHT_BROWSERS_PATH that deploy.sh uses when installing
# Chromium. Without this, Playwright defaults to ~/.cache/ms-playwright
# (i.e. /home/visual-compare/.cache/…) where nothing was installed, and
# capture jobs fail with "Executable doesn't exist".
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/visual-compare/.cache/ms-playwright

# ImageMagick 7 lives at /usr/local/bin/magick (a wrapper script — see
# cloud-init.api.yaml). The wrapper scopes LD_LIBRARY_PATH / MAGICK_* env
# to the magick invocation only, so the bundle's libs don't reach
# spawned Chromium (whose own glib has a conflicting ABI). The systemd
# unit deliberately does NOT export those vars — keep them out.

ExecStart=/usr/bin/node dist/index.js

Restart=on-failure
RestartSec=5s
# Playwright + ImageMagick can spike memory under load. Give the process
# headroom on the DEV1-M (4 GB) tier without letting a runaway eat the box.
MemoryMax=3G
# A failing comparison shouldn't take the API down — but a wedged process
# should be killable. 60 s is plenty for in-flight HTTP to drain.
TimeoutStopSec=60s

# Hardening — restrict filesystem and capabilities. /mnt/data is writable
# because that's where SQLite + the image tree + the lm-last-use file live.
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/mnt/data /tmp
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
